# sailingsa/deploy/sailingsa-admin-api.service
[Unit]
Description=SailingSA Admin API (port 8002; nginx /admin/)
After=network.target postgresql.service sailingsa-api.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/sailingsa/api
Environment="PATH=/var/www/sailingsa/api/venv/bin"
Environment="DB_URL=postgresql://sailors_user:SailSA_Pg_Beta2026@localhost:5432/sailors_master"
Environment="MAIN_API_URL=http://127.0.0.1:8000"
ExecStart=/var/www/sailingsa/api/venv/bin/uvicorn admin_api:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=1
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
